# .deploy/supervisord.ini
[supervisord]
logfile = /app/data/log/supervisord.log
logfile_maxbytes = 10MB
logfile_backups = 10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = true
childlogdir = /app/data/log/
environment = TZ=Asia/Shanghai

[inet_http_server]
port = 127.0.0.1:9005

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = http://127.0.0.1:9005

[program:gunicorn]
command = gunicorn SZUHelperBE.wsgi -c .deploy/gunicorn_config.py
directory = /app/
stdout_logfile = /app/data/log/gunicorn.log
stderr_logfile = /app/data/log/gunicorn.log
autostart = true
autorestart = true
startsecs = 5
stopwaitsecs = 5
killasgroup = true

[program:scheduler]
command = python manage.py runapscheduler
directory = /app/
stdout_logfile = /app/data/log/scheduler.log
stderr_logfile = /app/data/log/scheduler.log
autostart = true
autorestart = true
startsecs = 5
stopwaitsecs = 5
killasgroup = true

